# ref/core/nix-daemon.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /nix/var/nix/profiles/default/bin/nix-daemon --daemon
