# ref/quassel/quassel.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /nix/var/nix/profiles/quassel/bin/quasselcore
